# board/bitbox/rootfs_overlay/etc/profile.d/startup.sh
sudo chmod +x /usr/bin/static-get
sudo chmod +x /usr/bin/pfetch
sudo udhcpc -i enp0s5