# board/bitbox/rootfs_overlay/etc/profile.d/startup.sh
#!/bin/ash -ex

sudo chmod +x /usr/bin/static-get
sudo chmod +x /usr/bin/pfetch
sudo udhcpc